# mysql/supervisord-sshd.conf
[program:sshd]
command=/supervisord-sshd.sh
numprocs=1
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autostart=true
autorestart=true